# crontab.cron
@reboot /home/pi/Documents/rtcTest/clock.sh
